# docker/offload/iceberg-offload-scheduler.service
[Unit]
Description=K2 Iceberg Offload Scheduler
Documentation=file:///home/rjdscott/Documents/projects/k2-market-data-platform/docs/operations/prefect-schedule-config.md
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=rjdscott
WorkingDirectory=/home/rjdscott/Documents/projects/k2-market-data-platform
ExecStart=/usr/bin/python3 /home/rjdscott/Documents/projects/k2-market-data-platform/docker/offload/scheduler.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iceberg-offload

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
